# weatherdash.service
[Unit]
Description=MyWeatherDash - Real-time Weather Dashboard
After=network.target mariadb.service
Wants=mariadb.service

[Service]
Type=simple
User=weatherdash
Group=weatherdash
WorkingDirectory=/home/weatherdash/weatherdash
ExecStart=/home/weatherdash/weatherdash/weatherdash
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weatherdash

# Health check (systemd 250+)
# If your systemd supports it, uncomment:
# WatchdogSec=30s

[Install]
WantedBy=multi-user.target
